# debian/aim-aid.conf
description "ACI Inconsistency Detector"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 0 0
post-stop exec sleep 5

chdir /var/run

exec start-stop-daemon --start \
  --exec /usr/bin/aim-aid -- \
  --config-file=/etc/aim/aim.conf \
  --log-file=/var/log/aim-aid.log
